# consul_client/files/consul-client.service
[Unit]
Requires=network-online.target
After=network-online.target
ConditionFileNotEmpty=/opt/consul/etc/client/configuration.hcl

[Service]
User=consul
Group=consul
ExecStart=/opt/consul/bin/consul agent -config-dir=/opt/consul/etc/client/
ExecReload=/usr/local/bin/consul reload
KillMode=process
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target